# playbooks/roles/solana_validator_bootstrap/templates/solana-validator.service.j2
[Unit]
Description=Solana {{ cluster_environment }} node
After=network.target syslog.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User={{ solana_user }}
LimitNOFILE=1024000
Environment="PATH=/bin:/usr/bin:{{ env_path }}"
Environment="SOLANA_METRICS_CONFIG=host={{ solana_metrics_url }}"
{% if jito_enable %}
Environment="BLOCK_ENGINE_URL={{ jito_block_engine_url }}"
Environment="RELAYER_URL={{ jito_relayer_url }}"
Environment="SHRED_RECEIVER_ADDR={{ jito_receiver_addr }}"
{% endif %}
ExecStart={{ env_path }}/solana-validator \
--identity {{ secrets_path }}/validator-keypair.json \
--vote-account {{ secrets_path }}/vote-account-keypair.json \
--authorized-voter /home/solana/.secrets/funded-validator-keypair.json \
--rpc-port {{ solana_rpc_port }} \
--only-known-rpc \
--expected-genesis-hash {{ expected_genesis_hash }} \
{% for entrypoint in entrypoints %}
--entrypoint {{ entrypoint }} \
{% endfor %}
{% for known_validator in known_validators %}
--known-validator {{ known_validator }} \
{% endfor %}
--wal-recovery-mode skip_any_corrupted_record \
--limit-ledger-size \
--log {{ validator_log_path }}/solana-validator.log \
--maximum-incremental-snapshots-to-retain 2 \
--ledger {{ ledger_path }} \
--accounts {{ accounts_path }} \
--snapshots {{ snapshots_path }} \
--dynamic-port-range {{ open_solana_ports_start }}-{{ open_solana_ports_end }} \
{% if jito_enable %}
--tip-payment-program-pubkey {{ jito_tip_payment_program_pubkey }} \
--tip-distribution-program-pubkey {{ jito_distribution_program_pubkey }} \
--merkle-root-upload-authority {{ jito_merkle_root_upload_authority }} \
--commission-bps {{ jito_commission_bps }} \
--relayer-url ${RELAYER_URL} \
--block-engine-url ${BLOCK_ENGINE_URL} \
--shred-receiver-address ${SHRED_RECEIVER_ADDR} \
{% endif %}
--expected-shred-version {{ expected_shred_version }}

[Install]
WantedBy=multi-user.target
